# mimi-router.service
# mimi-router.service — mimi-router API Gateway
# File path: /etc/systemd/system/mimi-router.service
# sudo systemctl daemon-reload
# sudo systemctl start mimi-router
# sudo systemctl enable mimi-router
# sudo systemctl status mimi-router
[Unit]
Description=mimi-router Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/mimi-router  # 注意修改路径
ExecStart=/path/to/mimi-router/mimi-router --port 3000 --log-dir /path/to/mimi-router/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
